# bp-phone-tools/etc/cron.d/bp-ad-asterisk-sync
* * * * *	root	/usr/bin/bp-ad-asterisk-sync -p -f
0 0 * * *	root	/usr/bin/bp-ad-asterisk-sync -v -V
