# bp-server-ldap-auth/etc/cron.d/refresh_group_cache
0 *	* * *	root	/opt/pbis/bin/ad-cache --delete-all > /dev/null; getent group | grep uni > /dev/null 2>&1
