# bp-backupmanager/etc/cron.d/bp-backupmanager
0 * * * *	root	/usr/bin/bp-backupmanager -r

